# scripts/trt-web-dashboard.service
[Unit]
Description=TRT Web Dashboard Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=joshuag
WorkingDirectory=/home/joshuag/Time-Resolution-Theory-Live-Proof
ExecStart=/usr/bin/python3 /home/joshuag/Time-Resolution-Theory-Live-Proof/scripts/web_server.py
Restart=always
RestartSec=10
StandardOutput=append:/home/joshuag/Time-Resolution-Theory-Live-Proof/scripts/web_server.log
StandardError=append:/home/joshuag/Time-Resolution-Theory-Live-Proof/scripts/web_server.log

[Install]
WantedBy=multi-user.target
